# systemd/lan-orangutan.service
[Unit]
Description=LAN Orangutan Network Discovery
Documentation=https://github.com/291-Group/LAN-Orangutan
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as root so nmap can read the ARP table, which is what supplies MAC
# addresses and manufacturer names. Without it you get IP addresses only.
User=root
Group=root
ExecStart=/usr/local/bin/orangutan serve
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lan-orangutan /etc/lan-orangutan
PrivateTmp=true

[Install]
WantedBy=multi-user.target
